# config/supervisord.onprem.conf
[supervisord]
nodaemon=true
user=opc200
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:gateway]
command=openclaw gateway start --config /etc/opc200/gateway.yml
autostart=true
autorestart=true
stdout_logfile=/var/log/opc200/gateway.log
stderr_logfile=/var/log/opc200/gateway.error.log
environment=HOME="/home/opc200",USER="opc200"

[program:journal-api]
command=/opt/opc200-venv/bin/python /app/src/api/journal_api.py --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
stdout_logfile=/var/log/opc200/journal-api.log
stderr_logfile=/var/log/opc200/journal-api.error.log
environment=PYTHONPATH="/app/src",OPC200_DATA_DIR="/data",QDRANT_HOST="%(ENV_QDRANT_HOST)s"

[program:metrics-server]
command=/opt/opc200-venv/bin/python /app/src/api/metrics_server.py --host 0.0.0.0 --port 9091 --with-app-metrics
autostart=true
autorestart=true
stdout_logfile=/var/log/opc200/metrics.log
stderr_logfile=/var/log/opc200/metrics.error.log
environment=PYTHONPATH="/app/src",OPC200_DATA_DIR="/data"
